# infra/systemd_backend.service
[Unit]
Description=Admin Portal Backend
After=network.target
[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/admin-backend
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
Environment=NODE_ENV=production
Environment=DB_HOST=127.0.0.1
Environment=DB_USER=root
Environment=DB_PASS=changeme
[Install]
WantedBy=multi-user.target
